# deploy/systemd/opportunity-scout-ui.service
[Unit]
Description=Opportunity Scout Streamlit UI
After=network.target opportunity-scout-api.service
Wants=opportunity-scout-api.service

[Service]
Type=simple
User=oppscout
Group=oppscout
WorkingDirectory=/opt/opportunity-scout
EnvironmentFile=/etc/opportunity-scout/opportunity-scout.env
ExecStart=/opt/opportunity-scout/.venv/bin/streamlit run ui/app.py --server.headless true --server.address 127.0.0.1 --server.port 8500
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
